# builds/amd64-1.1/usr/lib/systemd/system/update-jenkins.timer
[Unit]
Description=Timer update Jenkins
Requires=update-jenkins.service

[Timer]
Unit=update-jenkins.service
OnCalendar=Sat *-*-1..8 02:00

[Install]
WantedBy=timers.target